# overlay/system/etc/init/fingerprint.conf
# fingerprint.conf - needed for FOD to work
start on filesystem
exec echo "fod_enable,1" > /sys/class/sec/tsp/cmd
